# halle-bot.service
[Unit]
Description=Halle Appointment Bot for Telegram
After=network.target
Documentation=https://github.com/yourusername/halle_bot/blob/main/README.md

[Service]
Type=simple
User=%u
WorkingDirectory=%h/halle_bot
Environment="PATH=%h/halle_bot/venv/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=%h/halle_bot/venv/bin/python %h/halle_bot/bot.py
Restart=on-failure
RestartSec=60
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
